# conf/crontab
* * * * * cd /var/www/html && php artisan schedule:run >> /dev/nul 2>&1